# cron.d/notify-stock
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
MAILTO=""
*/1	*	*	*	*	python /usr/local/notifyapp/run.py >> /var/log/cron.log 2>&1
